# dockerfs/etc/services.d/epson_projector/run
#!/command/execlineb -P

python /usr/src/app/hdfury_vrroom_mqtt.py
