# dockerfs/etc/services.d/hdfury_bridge/run
#!/command/with-contenv /bin/sh

python /usr/src/app/hdfury_vrroom_mqtt.py
